# config/country.conf
start on runlevel [2345]
stop on runlevel [016]
respawn

env PORT=5000
env NODE_ENV=production

setuid deployer

chdir /home/deployer/app
exec npm start
